# node.conf
# start apache on vagrant mounted
start on vagrant-ready
exec sudo /home/vagrant/node-v0.10.30-linux-x64/bin/node /vagrant/server.js
